# systemd/gate-tracker-watcher.service
[Unit]
Description=Gate Tracker GPIO Watcher
After=network-online.target gate-tracker-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/Gate-Tracker-Kinrise/src
Environment=DEVICE=entrance
ExecStart=/usr/bin/python3 /opt/Gate-Tracker-Kinrise/src/gpio_watcher.py --device ${DEVICE}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
